# tests/load/gcp/loki-log/locust-loggers.tpl.target
[Unit]
Description=A collection of locust logging sources
Wants=%{ for i in range(NUM_USERS) ~}locust@${i}.service %{ endfor ~}
