# tests/load/gcp/prom-scrape/avalanche-targets.tpl.target
[Unit]
Description=A collection of avalanche scrape targets
Wants=%{ for i in range(NUM_TARGETS) ~}avalanche@${9001 + i}.service %{ endfor ~}
